# ubuntu24/rootfs-overlay/etc/systemd/system/ModemManager.service.d/equuleus.conf
[Unit]
Wants=equuleus-modem.service
After=equuleus-modem.service
