# rootfs-overlay/etc/systemd/system/rmtfs.service.d/equuleus.conf
[Unit]
Before=equuleus-mss.service

[Service]
ExecStart=
# Never pass -s: stopping rmtfs must not stop the running MSS remoteproc.
ExecStart=/usr/bin/rmtfs -r -P

